# linux/udev/00-usb-keyboard.rules
ACTION=="add", ATTRS{idVendor}=="****", ATTRS{idProduct}=="****", RUN+="/home/jasalt/dotfiles/linux/udev/kbd_udev", OWNER="jasalt"
