# linux/udev/85-my_usb_rule.rules
ACTION=="add", SUBSYSTEM=="usb", SYSFS{idVendor}=="046d", SYSFS{idProduct}=="c52b", RUN+="/home/jarkko/dotfiles/linux/input-config.sh"
